# rootfs-overlay/etc/profile.d/profile.sh
alias ll='ls --color=tty -laFh'
alias ls='ls --color=tty -F'
